# roles/os/swap/disable/files/disable-swap.service
[Unit]
Description=Disable swap during boot

[Service]
Type=oneshot
ExecStart=/usr/sbin/swapoff --all

[Install]
WantedBy=multi-user.target
